# trunk/sources/cron/rc.cron
#!/bin/sh
#
# Startup script for cron
#
# Stop myself if running
if [ -n "`pidof cron`" ]; then
    /bin/killall cron 2>/dev/null
fi

/opt/sbin/cron
