# trunk/sources/ipac-ng/crontab
0,5,10,15,20,25,30,35,40,45,50,55 * * * * root /opt/sbin/fetchipac > /dev/null 2>&1
