# sources/cron/rc.cron
#!/bin/sh
#
# Startup script for cron
#
# Stop myself if running
PIDFILE=%OPTWARE_TARGET_PREFIX%/var/run/cron.pid
[ -f ${PIDFILE} ] && kill `cat ${PIDFILE}`
%OPTWARE_TARGET_PREFIX%/sbin/cron
